# supervisord.conf
[supervisord]
nodaemon=true

[program:shell]
command=/bin/bash

# Add additional program configurations here as needed
